# duplicacy-prune.timer
[Unit]
Description=Run duplicacy-prune on timer

[Timer]
OnBootSec=15min
OnUnitActiveSec=1day

[Install]
WantedBy=timers.target
