# duplicacy-backup.service
[Unit]
Description=run duplicacy backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/duplicacy-backup.sh
