# duplicacy-backup.timer
[Unit]
Description=Run duplicacy-backup on timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
