# duplicacy-prune.service
[Unit]
Description=run duplicacy prune

[Service]
Type=oneshot
ExecStart=/usr/local/bin/duplicacy-prune.sh
